# service/autossh.service
[Unit]
Description=rclone-over-autossh tunnel service
After=network.target

[Service]
EnvironmentFile=/etc/rclone-over-autossh/service.conf

# 基础命令配置
ExecStart=/usr/bin/autossh -M 0 \
    -o "ServerAliveInterval 30" \
    -o "ServerAliveCountMax 3" \
    -o "ExitOnForwardFailure yes" \
    -o "StrictHostKeyChecking no" \
    -o "UserKnownHosts /dev/null" \
    -NT \
    -i /etc/rclone-over-autossh/ssh_key \
    -R $RemoteAddress:$LocalAddress \
    $RemoteHostSpec

# 自动重启配置
Restart=always
RestartSec=10
StartLimitInterval=0
StartLimitBurst=0

# 安全加固
User=shared
Group=shared
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/log/rclone-over-autossh
ReadOnlyPaths=/etc/ssh
ReadOnlyPaths=/etc/rclone-over-autossh

# 日志输出
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rclone-over-autossh

# 资源限制
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target